# mod/hakchi/rootfs/etc/udev/rules.d/90-network.rules
ACTION=="add|remove", SUBSYSTEM=="net", RUN+="/bin/service network udev"
